# home/.config/autostart/xcape.desktop
[Desktop Entry]
Name=XCape
Exec=/bin/xcape
Terminal=false
Type=Application
StartupNotify=false
